# .config/systemd/user/config-backup.service
[Unit]
Description=Back up .config to S3 with restic
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/home/po/bin/backup-config-to-s3
