# .config/systemd/user/send-to-kindle.service
[Unit]
Description=Send content to Kindle

[Service]
Type=oneshot
ExecStart=/home/po/bin/send-to-kindle
