# .config/systemd/user/send-economist-to-kindle.service
[Unit]
Description=Send The Economist to Kindle

[Service]
Type=oneshot
ExecStart=%h/bin/send-economist-to-kindle
